# telegram-http-server.service
[Unit]
Description=Telegram File Download Bot HTTP Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/telegram-file-bot
Environment="PATH=/root/telegram-file-bot/venv/bin"
ExecStart=/root/telegram-file-bot/venv/bin/python http_server.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-http-server

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target