# telegram-file-bot.service
[Unit]
Description=Telegram File Download Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/telegram-file-bot
Environment="PATH=/root/telegram-file-bot/venv/bin"
ExecStart=/root/telegram-file-bot/venv/bin/python bot.py
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-file-bot

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target